# root/etc/s6-overlay/s6-rc.d/init-privoxy-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

set -euo pipefail

# shellcheck source=/scripts/logging.sh
source /scripts/logging.sh

if [ "${VPN_ENABLED:-false}" != "true" ] || [ "${PRIVOXY_ENABLED:-false}" != "true" ]; then
    debug_log "[Privoxy] disabled"
    exit 0
fi

config_file="/etc/privoxy/config"

log "[Privoxy] OpenVPN enabled, setting up Privoxy config"

# Bind on eth0 IPv4 address if available, otherwise 0.0.0.0
adr="$(ip -4 a show eth0 2>/dev/null | awk '/inet /{print $2}' | head -n1 | cut -d/ -f1 || true)"
listen_ip="${adr:-0.0.0.0}"
listen_port="${PRIVOXY_PORT:-8118}"

# Set the port for the IPv4 interface
sed -i -E "s/^listen-address\s+127\.0\.0\.1:.*/listen-address ${listen_ip}:${listen_port}/" "${config_file}" || true

# Remove IPv6 listen (keep it simple)
sed -i -E "/^listen-address\s+\[::1\].*/d" "${config_file}" || true

# Make sure we don't daemonize
if ! grep -q '^toggle' "${config_file}"; then
    debug_log "[Privoxy] config file format unexpected (no toggle line)"
fi

# Allow LAN access by default
if grep -q '^permit-access' "${config_file}"; then
    sed -i -E 's/^permit-access.*/permit-access 0.0.0.0\/0/' "${config_file}" || true
else
    echo 'permit-access 0.0.0.0/0' >> "${config_file}"
fi

debug_log "[Privoxy] listen-address=${listen_ip}:${listen_port}"
