# root/etc/s6-overlay/s6-rc.d/init-openvpn-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

set -euo pipefail

# shellcheck source=/scripts/logging.sh
source /scripts/logging.sh

vpn_dir="${VPN_DIR:-/app/config/openvpn}"
mkdir -p "${vpn_dir}"
chown -R tuliprox:tuliprox "${vpn_dir}" 2>/dev/null || true

if [ "${VPN_ENABLED:-false}" != "true" ]; then
    log "[OpenVPN] VPN disabled"
    exit 0
fi

vpn_provider="${OPENVPN_PROVIDER:-CUSTOM}"
vpn_provider_lc="$(echo "${vpn_provider}" | tr '[:upper:]' '[:lower:]')"

# Ensure tun exists
if [ ! -c /dev/net/tun ]; then
    mkdir -p /dev/net
    mknod /dev/net/tun c 10 200
    chmod 0666 /dev/net/tun
fi

# DNS override
if [ -n "${NAME_SERVERS:-}" ]; then
    : > /etc/resolv.conf
    IFS=',' read -ra servers <<< "${NAME_SERVERS}"
    for s in "${servers[@]}"; do
        s="$(echo "${s}" | sed -e 's~^[ \t]*~~;s~[ \t]*$~~')"
        [ -n "${s}" ] || continue
        echo "nameserver ${s}" >> /etc/resolv.conf
    done
fi

if [ "${vpn_provider_lc}" != "custom" ]; then
    provider_script="/etc/openvpn/${vpn_provider_lc}/update.sh"
    provider_dest="${vpn_dir}/${vpn_provider_lc}"

    # Allow forced re-download
    if [ "${OPENVPN_FORCE_UPDATE:-false}" = "true" ] && [ -d "${provider_dest}" ]; then
        log "[OpenVPN] OPENVPN_FORCE_UPDATE=true, removing cached configs"
        rm -rf "${provider_dest}"
    fi

    provider_tmp="${provider_dest}.tmp"

    if [ ! -d "${provider_dest}" ]; then
        # If a completed .tmp dir exists (download finished but mv didn't), just rename it
        if [ -f "${provider_tmp}/list.txt" ]; then
            log "[OpenVPN] completing previous download for ${vpn_provider_lc}"
            chown -R tuliprox:tuliprox "${provider_tmp}" 2>/dev/null || true
            mv "${provider_tmp}" "${provider_dest}"
        else
            if [ ! -x "${provider_script}" ]; then
                log "[OpenVPN] no update script for provider: ${vpn_provider_lc}"
                exit 1
            fi

            # Download directly into a tmp dir on the volume, then atomically rename
            # so a partial/interrupted download never blocks future re-downloads
            rm -rf "${provider_tmp}"
            mkdir -p "${provider_tmp}"

            # Copy provider scripts into the tmp dir so update.sh runs with volume as CWD
            cp /etc/openvpn/${vpn_provider_lc}/*.sh "${provider_tmp}/" 2>/dev/null || true
            chmod +x "${provider_tmp}"/*.sh 2>/dev/null || true

            log "[OpenVPN] downloading ${vpn_provider_lc} configs to ${provider_dest}"
            (cd "${provider_tmp}" && bash "${provider_script}")

            chown -R tuliprox:tuliprox "${provider_tmp}" 2>/dev/null || true
            mv "${provider_tmp}" "${provider_dest}"
        fi
        log "[OpenVPN] ${vpn_provider_lc} configs ready at ${provider_dest}"
    else
        debug_log "[OpenVPN] using cached ${vpn_provider_lc} configs from ${provider_dest}"
    fi

    if [ -d "${provider_dest}" ]; then
        : > "${provider_dest}/list.txt"
        find "${provider_dest}" -maxdepth 2 -type f -name '*.ovpn' | while read -r cfg; do
            echo "${cfg#${provider_dest}/}" >> "${provider_dest}/list.txt"
        done
    fi
else
    # remove macOS metadata files if present
    rm -f "${vpn_dir}/._"* 2>/dev/null || true
fi

# Create credentials file(s) if env vars are supplied
if [ -n "${OPENVPN_USERNAME:-}" ] && [ -n "${OPENVPN_PASSWORD:-}" ]; then
    cred_file="${vpn_dir}/${vpn_provider_lc}-openvpn-credentials.txt"
    if [ "${vpn_provider_lc}" = "custom" ]; then
        cred_file="${vpn_dir}/custom-openvpn-credentials.txt"
    fi

    printf '%s\n%s\n' "${OPENVPN_USERNAME}" "${OPENVPN_PASSWORD}" > "${cred_file}"
    chmod 600 "${cred_file}"
    chown tuliprox:tuliprox "${cred_file}" 2>/dev/null || true

    if echo "${OPENVPN_USERNAME}" | grep -Eq '[^a-zA-Z0-9@]+'; then
        debug_log "[OpenVPN] username contains special characters which could cause auth issues"
    fi
    if echo "${OPENVPN_PASSWORD}" | grep -Eq '[^a-zA-Z0-9@]+'; then
        debug_log "[OpenVPN] password contains special characters which could cause auth issues"
    fi
fi

# Local routes (so you can still reach LAN services when VPN is up)
if [ -n "${LOCAL_NETWORK:-}" ]; then
    # Explicitly exclude tun0 so we always route via the physical interface (eth0)
    gw="$(ip route list match 0.0.0.0/0 | awk '$5 != "tun0" {print $3; exit}' || true)"
    intf="$(ip route list match 0.0.0.0/0 | awk '$5 != "tun0" {print $5; exit}' || true)"

    if [ -n "${gw}" ] && [ -n "${intf}" ]; then
        for net in ${LOCAL_NETWORK//,/ }; do
            log "[OpenVPN] adding route ${net} via ${gw} dev ${intf}"
            ip route add "${net}" via "${gw}" dev "${intf}" 2>/dev/null || true
        done
    else
        log "[OpenVPN] WARNING: could not determine gateway for LOCAL_NETWORK routes (tun0 excluded)"
    fi
fi

# Validate config + credentials
validation_out="$(/scripts/openvpn-config-validation.sh)"
validation_exit=$?

if [ ${validation_exit} -ne 0 ]; then
    echo "${validation_out}"
    exit ${validation_exit}
fi

# shellcheck disable=SC1090
eval "${validation_out}"

if [ -z "${OPENVPN_CONFIG_FILE:-}" ]; then
    log "[OpenVPN] validation did not return OPENVPN_CONFIG_FILE"
    exit 1
fi

log "[OpenVPN] using config: ${OPENVPN_CONFIG_FILE}"

# Clean only the selected config (not all provider configs)
debug_log "[OpenVPN] cleaning selected config"
/scripts/openvpn-config-clean.sh "${OPENVPN_CONFIG_FILE}" || true

# Provider-specific fixes on the selected config only
cred_file="${vpn_dir}/${vpn_provider_lc}-openvpn-credentials.txt"
case "${vpn_provider_lc}" in
    nordvpn)
        if grep -q '^cipher AES-256-CBC' "${OPENVPN_CONFIG_FILE}" 2>/dev/null; then
            sed -i 's/^cipher AES-256-CBC$/cipher AES-256-GCM\ndata-ciphers AES-256-GCM/' "${OPENVPN_CONFIG_FILE}" || true
        fi
        sed -i "s|^auth-user-pass.*|auth-user-pass ${cred_file}|" "${OPENVPN_CONFIG_FILE}" || true
        ;;
    surfshark)
        sed -i "s/AES-256-CBC/AES-128-GCM/g" "${OPENVPN_CONFIG_FILE}" || true
        sed -i "s|auth-user-pass.*|auth-user-pass ${cred_file}|g" "${OPENVPN_CONFIG_FILE}" || true
        ;;
    pia)
        sed -i "s|auth-user-pass.*|auth-user-pass ${cred_file}|g" "${OPENVPN_CONFIG_FILE}" || true
        sed -i "s|ca ca\.rsa\.\([0-9]*\)\.crt|ca ${vpn_dir}/pia/ca.rsa.\1.crt|g" "${OPENVPN_CONFIG_FILE}" || true
        sed -i "s|crl-verify crl\.rsa\.\([0-9]*\)\.pem|crl-verify ${vpn_dir}/pia/crl.rsa.\1.pem|g" "${OPENVPN_CONFIG_FILE}" || true
        ;;
    ipvanish)
        sed -i "s|auth-user-pass.*|auth-user-pass ${cred_file}|g" "${OPENVPN_CONFIG_FILE}" || true
        ;;
    vyprvpn)
        sed -i "s|auth-user-pass.*|auth-user-pass ${cred_file}|g" "${OPENVPN_CONFIG_FILE}" || true
        ;;
    protonvpn)
        sed -i "s|auth-user-pass.*|auth-user-pass ${cred_file}|g" "${OPENVPN_CONFIG_FILE}" || true
        ;;
esac

# Persist resolved paths to s6 container environment so longruns inherit them
mkdir -p /run/s6/container_environment
printf '%s' "${OPENVPN_CONFIG_FILE}" > /run/s6/container_environment/OPENVPN_CONFIG_FILE
if [ -n "${OPENVPN_AUTH_FILE:-}" ]; then
    printf '%s' "${OPENVPN_AUTH_FILE}" > /run/s6/container_environment/OPENVPN_AUTH_FILE
fi
if [ -n "${LOCAL_NETWORK:-}" ]; then
    printf '%s' "${LOCAL_NETWORK}" > /run/s6/container_environment/LOCAL_NETWORK
fi
