# root/etc/s6-overlay/s6-rc.d/svc-privoxy/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

set -euo pipefail

# shellcheck source=/scripts/logging.sh
source /scripts/logging.sh

if [ "${VPN_ENABLED:-false}" != "true" ] || [ "${PRIVOXY_ENABLED:-false}" != "true" ]; then
    debug_log "[Privoxy] disabled (sleeping)"
    exec sleep infinity
fi

timeout_secs="${PRIVOXY_STARTUP_DELAY_SECS:-30}"
log "[Privoxy] waiting up to ${timeout_secs}s for tun interface before starting"
elapsed=0
while [ ${elapsed} -lt ${timeout_secs} ]; do
    if ip link show tun0 >/dev/null 2>&1; then
        log "[Privoxy] tun0 is up after ${elapsed}s"
        break
    fi
    sleep 1
    elapsed=$(( elapsed + 1 ))
done
if ! ip link show tun0 >/dev/null 2>&1; then
    log "[Privoxy] WARNING: tun0 not found after ${timeout_secs}s, starting anyway"
fi

port="${PRIVOXY_PORT:-8118}"
log "[Privoxy] starting on port ${port}"

log_dir="${PRIVOXY_LOG_DIR:-/app/config/logs/privoxy}"
mkdir -p "${log_dir}"
chown -R tuliprox:tuliprox "${log_dir}" 2>/dev/null || true

exec /usr/sbin/privoxy \
    --pidfile /var/run/privoxy.pid \
    --no-daemon \
    /etc/privoxy/config 2>&1 | s6-log -b n20 s1000000 T "${log_dir}"
