# root/etc/s6-overlay/s6-rc.d/svc-openvpn/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
#
# Runtime VPN Config Switching Support
#
# This modification allows changing VPN config at runtime by:
# 1. Writing the desired config path to ${VPN_DIR}/.switch-config
# 2. Sending SIGHUP or restarting the openvpn service
#
# The script checks for the override file on each restart, enabling
# hot-swapping VPN servers without container restart.

set -euo pipefail

# shellcheck source=/scripts/logging.sh
source /scripts/logging.sh

if [ "${VPN_ENABLED:-false}" != "true" ]; then
    log "[OpenVPN] VPN disabled (sleeping)"
    exec sleep infinity
fi

# Check for runtime config override
# This file can be written by external tools to change VPN config dynamically
CONFIG_OVERRIDE_FILE="${VPN_DIR:-/data/config/openvpn}/.switch-config"
if [ -f "$CONFIG_OVERRIDE_FILE" ]; then
    OVERRIDE_CONFIG=$(cat "$CONFIG_OVERRIDE_FILE" 2>/dev/null | tr -d '\n')
    if [ -n "$OVERRIDE_CONFIG" ] && [ -f "$OVERRIDE_CONFIG" ]; then
        log "[OpenVPN] using runtime override config: $OVERRIDE_CONFIG"
        OPENVPN_CONFIG_FILE="$OVERRIDE_CONFIG"
        # Clear the override so it's only used once
        rm -f "$CONFIG_OVERRIDE_FILE"
    fi
fi

# Fall back to standard env var if no override
if [ -z "${OPENVPN_CONFIG_FILE:-}" ]; then
    log "[OpenVPN] OPENVPN_CONFIG_FILE not set — init-openvpn-config may have failed"
    exec sleep infinity
fi

openvpn_args=(
    --reneg-sec 0
    --disable-occ
    --auth-nocache
    --mute-replay-warnings
    --data-ciphers-fallback AES-128-GCM
    --pull-filter ignore 'route-ipv6'
    --pull-filter ignore 'ifconfig-ipv6'
    --pull-filter ignore 'tun-ipv6'
    --pull-filter ignore 'dhcp-option DNS6'
    --pull-filter ignore 'reneg-sec'
    --pull-filter ignore 'persist-tun'
    --pull-filter ignore 'block-outside-dns'
)

# Allow user custom options
if [ -n "${OPENVPN_OPTIONS:-}" ]; then
    # shellcheck disable=SC2206
    extra=( ${OPENVPN_OPTIONS} )
    openvpn_args+=( "${extra[@]}" )
fi

openvpn_args+=( --config "${OPENVPN_CONFIG_FILE}" )

# Inject --auth-user-pass if credentials file exists and config doesn't already have it
if [ -n "${OPENVPN_AUTH_FILE:-}" ] && [ -f "${OPENVPN_AUTH_FILE}" ]; then
    if ! grep -Eq '^auth-user-pass' "${OPENVPN_CONFIG_FILE}" 2>/dev/null; then
        openvpn_args+=( --auth-user-pass "${OPENVPN_AUTH_FILE}" )
        debug_log "[OpenVPN] injecting --auth-user-pass ${OPENVPN_AUTH_FILE}"
    fi
fi

debug_log "[OpenVPN] OPENVPN_CONFIG_FILE=${OPENVPN_CONFIG_FILE}"
log "[OpenVPN] starting openvpn"

log_dir="${VPN_LOG_DIR:-/app/config/logs/openvpn}"
mkdir -p "${log_dir}"
chown -R tuliprox:tuliprox "${log_dir}" 2>/dev/null || true

exec openvpn "${openvpn_args[@]}" 2>&1 | s6-log -b n20 s1000000 T "${log_dir}"
