# root/etc/s6-overlay/s6-rc.d/svc-openvpn/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

set -euo pipefail

# shellcheck source=/scripts/logging.sh
source /scripts/logging.sh

if [ "${VPN_ENABLED:-false}" != "true" ]; then
    log "[OpenVPN] VPN disabled (sleeping)"
    exec sleep infinity
fi

# OPENVPN_CONFIG_FILE is written to the s6 container environment by init-openvpn-config
if [ -z "${OPENVPN_CONFIG_FILE:-}" ]; then
    log "[OpenVPN] OPENVPN_CONFIG_FILE not set — init-openvpn-config may have failed"
    exec sleep infinity
fi

openvpn_args=(
    --reneg-sec 0
    --disable-occ
    --auth-nocache
    --mute-replay-warnings
    --data-ciphers-fallback AES-128-GCM
    --pull-filter ignore 'route-ipv6'
    --pull-filter ignore 'ifconfig-ipv6'
    --pull-filter ignore 'tun-ipv6'
    --pull-filter ignore 'dhcp-option DNS6'
    --pull-filter ignore 'reneg-sec'
    --pull-filter ignore 'persist-tun'
    --pull-filter ignore 'block-outside-dns'
    --pull-filter ignore 'redirect-gateway'
)

# Allow user custom options
if [ -n "${OPENVPN_OPTIONS:-}" ]; then
    # shellcheck disable=SC2206
    extra=( ${OPENVPN_OPTIONS} )
    openvpn_args+=( "${extra[@]}" )
fi

openvpn_args+=( --config "${OPENVPN_CONFIG_FILE}" )

# Inject --auth-user-pass if credentials file exists and config doesn't already have it
if [ -n "${OPENVPN_AUTH_FILE:-}" ] && [ -f "${OPENVPN_AUTH_FILE}" ]; then
    if ! grep -Eq '^auth-user-pass' "${OPENVPN_CONFIG_FILE}" 2>/dev/null; then
        openvpn_args+=( --auth-user-pass "${OPENVPN_AUTH_FILE}" )
        debug_log "[OpenVPN] injecting --auth-user-pass ${OPENVPN_AUTH_FILE}"
    fi
fi

debug_log "[OpenVPN] OPENVPN_CONFIG_FILE=${OPENVPN_CONFIG_FILE}"
log "[OpenVPN] starting openvpn"

log_dir="${VPN_LOG_DIR:-/app/config/logs/openvpn}"
mkdir -p "${log_dir}"
chown -R tuliprox:tuliprox "${log_dir}" 2>/dev/null || true

exec openvpn "${openvpn_args[@]}" 2>&1 | s6-log -b n20 s1000000 T "${log_dir}"
